# systemd/etcd.service
[Unit]
Description=etcd - distributed key-value store

[Service]
ExecStart=etcd -config-file /etc/etcd.conf.yml

[Install]
WantedBy=default.target